# airootfs/etc/systemd/system/apparmor.service.d/override.conf
[Unit]
After=local-fs.target

[Service]
ExecStartPre=/usr/sbin/apparmor_parser --skip-cache --write-cache /etc/apparmor.d
